# packages/applier/applier-0.1.1.tar.gz/applier-0.1.1/configs/applier.socket
[Unit]
Description=applier socket

[Socket]
ListenStream=/run/applier.sock
SocketUser=www-data
SocketMode=600

[Install]
WantedBy=sockets.target
