# packages/rssbot/rssbot-201.tar.gz/rssbot-201/files/rssbot.service
[Unit]
Description=RSSBOT - feeding rss into your irc channel
After=multi-user.target

[Service]
Type=simple
User=rssbot
Group=rssbot
DynamicUser=yes
StateDirectory=rssbot
LogsDirectory=rssbot
CacheDirectory=rssbot
ExecStart=/usr/local/bin/rssbot
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
